# etc/systemd/supermux.service
# supermux systemd unit — TEMPLATE.
#
# This is a template, not a ready-to-install unit. `scripts/deploy.sh`
# substitutes the placeholders (__SERVICE_USER__, __USER_HOME__, __DATA_DIR__,
# __PROTECT_HOME__, __READ_WRITE_PATHS__) from your environment (or .env) and
# installs the result to /etc/systemd/system/supermux.service on the host.
#
# HOME vs DATA_DIR — DELIBERATELY SEPARATE.
#   __USER_HOME__  is the service user's actual login home (/home/<user> or
#                  /root). The unit's WorkingDirectory and $HOME both point
#                  here so spawned shells/agents (claude, codex, tmux children)
#                  resolve their dotfiles (.bash_history, .claude/,
#                  .claude.json, .cache/, .local/, …) in the conventional
#                  place — NOT inside the supermux data dir.
#   __DATA_DIR__   is supermux's own data dir (SQLite, auth token, config.toml,
#                  uploads). The server reads SUPERMUX_DATA_DIR from the
#                  environment; it does not — and should not — share storage
#                  with the user's shell dotfiles. Keeping the two separate
#                  also means a `rm -rf $SUPERMUX_DATA_DIR` does not nuke the
#                  user's shell history or agent caches.
#
# SECURITY: by default the service is NOT run as root. supermux-server exposes
# a full filesystem read/write API (GET /api/ls, GET/PUT /api/file,
# POST /api/fs/upload, DELETE /api/fs/delete) and spawns child processes
# (tmux, git, claude, codex). Running as uid 0 would make any auth-token leak,
# path-jail escape or command-injection root-equivalent on the host. The
# service therefore runs as an unprivileged account, with systemd sandboxing
# locked down to exactly what the service needs:
#   - read+write only the configured paths — ProtectSystem=strict makes the
#     rest of the filesystem read-only, ReadWritePaths scopes the writes;
#   - no privilege escalation (NoNewPrivileges), private /tmp, kernel + cgroup
#     tunables protected, address families restricted to what an HTTP service
#     plus the tmux AF_UNIX socket need.
#
# ROOT DEPLOYS: if you explicitly opt in to running as root (deploy.sh requires
# SUPERMUX_ALLOW_ROOT=1), the deploy renders __PROTECT_HOME__=false (because
# ProtectHome=true would mask /root and the unit could not chdir to its
# WorkingDirectory) and __READ_WRITE_PATHS__=/root (so tmux/claude/git children
# can write in arbitrary subdirs under /root). All other hardening directives
# still apply. The remaining sandbox is still meaningful, but the privilege
# trade-off is on the operator.
[Unit]
Description=supermux (Rust)
After=network.target
# OPTIONAL: if you expose supermux via `tailscale serve`, uncomment the next
# two lines so the unit waits for Tailscale to be up. deploy.sh appends these
# automatically when SUPERMUX_USE_TAILSCALE=1.
# After=tailscaled.service
# Wants=tailscaled.service

[Service]
Type=simple
# Run as the unprivileged deployment account. This account owns the data dir
# and the tmux server it drives — no root privilege required.
User=__SERVICE_USER__
# WorkingDirectory + HOME point at the SERVICE USER's actual login home (NOT
# the supermux data dir). This keeps spawned shells/agents (claude, codex,
# tmux children) writing their dotfiles (.bash_history, .claude/,
# .claude.json, .cache/, .local/, …) into the conventional place instead of
# polluting the data dir. SUPERMUX_DATA_DIR below is what supermux-server
# actually reads for its own storage.
WorkingDirectory=__USER_HOME__
Environment=HOME=__USER_HOME__
Environment=SUPERMUX_DATA_DIR=__DATA_DIR__
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/usr/local/bin/supermux-server
Restart=always
RestartSec=2
LimitNOFILE=65536
StandardOutput=journal
StandardError=journal

# ── hardening ───────────────────────────────────────────────────────────────
# No process spawned by the service may gain privileges via setuid/setgid/caps.
NoNewPrivileges=true
# Whole filesystem read-only (incl. /usr, /boot, /etc); only the paths in
# ReadWritePaths below are writable. ProtectSystem=strict deliberately breaks
# the files API writing outside the data dir — that is the intended jail.
ProtectSystem=strict
# Home-dir protection. Rendered by deploy.sh per deploy mode:
#   - non-root → `tmpfs`: /home is replaced by an empty tmpfs so OTHER users'
#     homes are hidden, and the __BIND_PATHS__ line below re-mounts ONLY this
#     service user's own home read-write (its WorkingDirectory + data dir +
#     spawned-agent dotfiles live there). `true` would mask the service user's
#     OWN home too — systemd then can't chdir to WorkingDirectory (exit
#     200/CHDIR) and ReadWritePaths can't un-mask it.
#   - root (SUPERMUX_ALLOW_ROOT=1) → `false`: /root stays reachable; no bind.
ProtectHome=__PROTECT_HOME__
# Re-mount the service user's own home past ProtectHome=tmpfs (non-root). For
# root deploys deploy.sh deletes this line entirely (ProtectHome=false).
__BIND_PATHS__
# Paths the service may write to. Defaults to the data dir for non-root
# deploys; deploy.sh broadens it for root (to /root) and exposes
# SUPERMUX_READ_WRITE_PATHS so operators can extend it for multi-project
# workflows (e.g. /home/supermux:/opt/projects). Multiple paths are
# whitespace-separated per systemd's ReadWritePaths= syntax.
ReadWritePaths=__READ_WRITE_PATHS__
# tmux keeps its server socket under $XDG_RUNTIME_DIR (/run/user/<uid>); the
# service must be able to create/connect to it to drive sessions.
RuntimeDirectory=supermux
# Private /tmp and /var/tmp namespace — no cross-service tmp snooping.
PrivateTmp=true
# Block writes to kernel tunables, kernel module load, and cgroup hierarchy.
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
# Localhost HTTP service + tmux AF_UNIX socket: AF_INET/AF_INET6 for the
# listener, AF_UNIX for tmux. Nothing else (no AF_NETLINK, AF_PACKET, …).
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
# No new namespaces, no realtime scheduling, no personality (ABI) changes,
# no setuid/setgid binaries, no executable+writable memory.
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
# Drop all Linux capabilities — the service binds an unprivileged port and
# needs none.
CapabilityBoundingSet=
AmbientCapabilities=
# Only the device nodes a normal userspace process needs (/dev/null, zero,
# random, …); no raw device access.
PrivateDevices=true
# Confine syscalls to the standard non-privileged set plus exec (needed to
# spawn tmux/git/claude/codex children).
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources @obsolete
SystemCallErrorNumber=EPERM
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
