# etc/systemd/supermux.service
# supermux systemd unit — TEMPLATE.
#
# This is a template, not a ready-to-install unit. `scripts/deploy.sh`
# substitutes the placeholders (__SERVICE_USER__, __USER_HOME__, __DATA_DIR__,
# __PROTECT_HOME__, __READ_WRITE_PATHS__) from your environment (or .env) and
# installs the result to /etc/systemd/system/supermux.service on the host.
#
# HOME vs DATA_DIR — DELIBERATELY SEPARATE.
#   __USER_HOME__  is the service user's actual login home (/home/<user> or
#                  /root). The unit's WorkingDirectory and $HOME both point
#                  here so spawned shells/agents (claude, codex, tmux children)
#                  resolve their dotfiles (.bash_history, .claude/,
#                  .claude.json, .cache/, .local/, …) in the conventional
#                  place — NOT inside the supermux data dir.
#   __DATA_DIR__   is supermux's own data dir (SQLite, auth token, config.toml,
#                  uploads). The server reads SUPERMUX_DATA_DIR from the
#                  environment; it does not — and should not — share storage
#                  with the user's shell dotfiles. Keeping the two separate
#                  also means a `rm -rf $SUPERMUX_DATA_DIR` does not nuke the
#                  user's shell history or agent caches.
#
# SECURITY: by default the service is NOT run as root. supermux-server exposes
# a full filesystem read/write API (GET /api/ls, GET/PUT /api/file,
# POST /api/fs/upload, DELETE /api/fs/delete) and spawns child processes
# (tmux, git, claude, codex). Running as uid 0 would make any auth-token leak,
# path-jail escape or command-injection root-equivalent on the host. The
# service therefore runs as an unprivileged account, with systemd sandboxing
# locked down to exactly what the service needs:
#   - read+write only the configured paths — ProtectSystem=strict makes the
#     rest of the filesystem read-only, ReadWritePaths scopes the writes;
#   - no privilege escalation (NoNewPrivileges), private /tmp, kernel + cgroup
#     tunables protected, address families restricted to what an HTTP service
#     plus the tmux AF_UNIX socket need.
#
# ROOT DEPLOYS: if you explicitly opt in to running as root (deploy.sh requires
# SUPERMUX_ALLOW_ROOT=1), the deploy renders __PROTECT_HOME__=false (because
# ProtectHome=true would mask /root and the unit could not chdir to its
# WorkingDirectory) and __READ_WRITE_PATHS__=/root (so tmux/claude/git children
# can write in arbitrary subdirs under /root). All other hardening directives
# still apply. The remaining sandbox is still meaningful, but the privilege
# trade-off is on the operator.
[Unit]
Description=supermux (Rust)
After=network.target
# OPTIONAL: if you expose supermux via `tailscale serve`, uncomment the next
# two lines so the unit waits for Tailscale to be up. deploy.sh appends these
# automatically when SUPERMUX_USE_TAILSCALE=1.
# After=tailscaled.service
# Wants=tailscaled.service

[Service]
Type=simple
# Run as the unprivileged deployment account. This account owns the data dir
# and the tmux server it drives — no root privilege required.
User=__SERVICE_USER__
# WorkingDirectory + HOME point at the SERVICE USER's actual login home (NOT
# the supermux data dir). This keeps spawned shells/agents (claude, codex,
# tmux children) writing their dotfiles (.bash_history, .claude/,
# .claude.json, .cache/, .local/, …) into the conventional place instead of
# polluting the data dir. SUPERMUX_DATA_DIR below is what supermux-server
# actually reads for its own storage.
WorkingDirectory=__USER_HOME__
Environment=HOME=__USER_HOME__
Environment=SUPERMUX_DATA_DIR=__DATA_DIR__
# Anchor tmux's control socket in the PERSISTENT data dir (not the default /tmp).
# With PrivateTmp=true below, /tmp is recreated every (re)start, which would
# orphan all running tmux sessions on each restart/deploy despite KillMode=process
# — the new instance could not reach the previous tmux server. A path under the
# bind-mounted data dir survives restarts, so sessions reconnect. (The server also
# defaults this in code; this line makes it explicit + covers any early tmux use.)
Environment=TMUX_TMPDIR=__DATA_DIR__/tmux
# The colon-separated project-dirs deploy.sh provisions (smart default
# <home>/projects; operators may set SUPERMUX_PROJECT_DIRS to a custom list,
# e.g. /opt/projects). The server reads this at runtime and exposes the FIRST
# entry to the SPA as `window._SUPERMUX_PROJECT_DIR` so Start-a-team can
# pre-fill its directory field with the projects root — the autocomplete then
# immediately surfaces the project repos for a one-click pick.
Environment=SUPERMUX_PROJECT_DIRS=__PROJECT_DIRS__
# ~/.local/bin first so an opt-in cargo-guard wrapper (deploy.sh installs it
# under SUPERMUX_INSTALL_CARGO_GUARD=1) is inherited by every spawned child,
# not just interactive bash that reads .bashrc.
Environment=PATH=__USER_HOME__/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
# Optional VAPID 'sub' contact (RFC-8292) for Web Push. Rendered by deploy.sh
# from SUPERMUX_PUSH_SUB (.env); if unset, deploy.sh DELETES the placeholder
# line below entirely (same delete-on-empty pattern used for the bind-paths
# directive) — an empty `Environment=` is a systemd parse error. The
# config.toml `push_sub` setting on the host remains the file-on-host alternative.
__PUSH_SUB_LINE__
ExecStart=/usr/local/bin/supermux-server
Restart=always
RestartSec=2
# KillMode=process — NOT the default control-group. supermux is a SESSION
# SUPERVISOR: it spawns a tmux server (+ claude/codex/shell panes) that must
# OUTLIVE a restart of supermux-server itself. The default control-group would
# SIGKILL the whole cgroup — tmux and every agent pane — on every deploy,
# crash-restart, or `systemctl restart`, so all sessions would come back
# `stopped`. With `process`, systemd signals ONLY the main supermux-server
# process; the tmux server keeps running and the freshly-started supermux
# re-attaches to it (the spawn path checks `tmux has-session` first). Trade-off:
# `systemctl stop` then leaves agents running (intended — that's the persistence
# guarantee). NOTE: a machine reboot still ends everything; surviving reboots is
# a separate restore-on-boot feature.
KillMode=process
LimitNOFILE=65536
# Bias the OOM-killer toward supermux's tree (children inherit). Under memory
# pressure the kernel prefers to evict a runaway agent/rustc here over sshd,
# keeping the box recoverable without a reboot.
OOMScoreAdjust=200
# Soft memory throttle for the whole cgroup (supermux + tmux + agents + their
# children). Rendered by deploy.sh from SUPERMUX_MEMORY_HIGH (default `infinity`
# = no cap, so the default behaviour is identical for every existing operator).
# When set (e.g. `6G` on an 8 GiB host, `12G` on a 16 GiB host), the kernel
# reclaims aggressively past the threshold but NEVER OOM-kills — so an
# over-tight value just slows things down, it doesn't break the service. Pairs
# with OOMScoreAdjust=200 above: this caps the cgroup's growth; OOMScoreAdjust
# decides who gets killed if the host as a whole still tips into OOM.
MemoryHigh=__MEMORY_HIGH__
# CPU de-prioritisation for the whole supermux tree (server + tmux + agents +
# their children). CPUWeight only bites UNDER CONTENTION: when the box is busy it
# proportionally favours everything else (sshd, the OS) over supermux's cgroup, so
# a runaway agent that pegs every core can't lock you out of the box or starve the
# server's own request handling — yet on an idle box it costs nothing (no cap).
# 50 = half the default 100 weight. This is the cheap, always-safe guard; it does
# NOT kill or hard-cap anything (a healthy busy agent still gets all spare CPU),
# unlike a CPUQuota ceiling. Static (no placeholder) so every deploy inherits it.
CPUWeight=50
StandardOutput=journal
StandardError=journal

# ── hardening ───────────────────────────────────────────────────────────────
# No process spawned by the service may gain privileges via setuid/setgid/caps.
NoNewPrivileges=true
# Whole filesystem read-only (incl. /usr, /boot, /etc); only the paths in
# ReadWritePaths below are writable. ProtectSystem=strict deliberately breaks
# the files API writing outside the data dir — that is the intended jail.
ProtectSystem=strict
# Home-dir protection. Rendered by deploy.sh per deploy mode:
#   - non-root → `tmpfs`: /home is replaced by an empty tmpfs so OTHER users'
#     homes are hidden, and the __BIND_PATHS__ line below re-mounts ONLY this
#     service user's own home read-write (its WorkingDirectory + data dir +
#     spawned-agent dotfiles live there). `true` would mask the service user's
#     OWN home too — systemd then can't chdir to WorkingDirectory (exit
#     200/CHDIR) and ReadWritePaths can't un-mask it.
#   - root (SUPERMUX_ALLOW_ROOT=1) → `false`: /root stays reachable; no bind.
ProtectHome=__PROTECT_HOME__
# Re-mount the service user's own home past ProtectHome=tmpfs (non-root). For
# root deploys deploy.sh deletes this line entirely (ProtectHome=false).
__BIND_PATHS__
# Paths the service may write to. Defaults to the data dir for non-root
# deploys; deploy.sh broadens it for root (to /root) and exposes
# SUPERMUX_READ_WRITE_PATHS so operators can extend it for multi-project
# workflows (e.g. /home/supermux:/opt/projects). Multiple paths are
# whitespace-separated per systemd's ReadWritePaths= syntax.
ReadWritePaths=__READ_WRITE_PATHS__
# /tmp + /var/tmp MUST be writable. With PrivateTmp=no (the default — see below),
# systemd no longer hands the service a private writable /tmp, and ProtectSystem=
# strict makes the real /tmp READ-ONLY unless it is listed here. The Claude Code
# agent harness writes its per-session/per-tool working state under /tmp
# (/tmp/claude-<uid>/...), so a read-only /tmp means claude launches but can never
# initialise — it hangs before painting its TUI and every tool call fails. This
# line restores a normal, persistent, writable /tmp (shared with the host; fine
# on a single-purpose box) so agents start, paint, and run tools, and that state
# survives a restart. (Under the hardened profile PrivateTmp=yes supplies its own
# writable /tmp, so this line is harmless there.)
ReadWritePaths=/tmp /var/tmp
# tmux keeps its server socket under $XDG_RUNTIME_DIR (/run/user/<uid>); the
# service must be able to create/connect to it to drive sessions.
RuntimeDirectory=supermux
# Private /tmp + /var/tmp namespace. DEFAULT OFF (__PRIVATE_TMP__, from
# SUPERMUX_HARDENED) and that is load-bearing, not a nicety: supermux runs a
# LONG-LIVED tmux server plus Claude Code agents, and the agent harness writes
# its per-tool working state under /tmp (/tmp/claude-<uid>/.../tasks). systemd
# tears down a PrivateTmp namespace on every service (re)start — so with it ON,
# each restart/deploy yanks /tmp out from under the SURVIVING tmux server and its
# agents: the agents can no longer mkdir in /tmp (every tool call fails ENOENT)
# and freshly-created panes never paint a TUI. Anchoring TMUX_TMPDIR in the data
# dir fixes the tmux socket but NOT the agents' /tmp. So PrivateTmp must stay off
# wherever sessions are meant to survive a restart; only the hardened profile
# (which assumes no long-lived agents across restarts) turns it on.
PrivateTmp=__PRIVATE_TMP__
# Block writes to kernel tunables, kernel module load, and cgroup hierarchy.
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
# Localhost HTTP service + tmux AF_UNIX socket: AF_INET/AF_INET6 for the
# listener, AF_UNIX for tmux. AF_NETLINK so agent dev tooling that enumerates
# network interfaces works — e.g. glibc getifaddrs()/`os.networkInterfaces()`,
# which Vite/Next/webpack-dev-server call to print the "Network:" URL. No
# AF_PACKET / raw families.
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX AF_NETLINK
# No new namespaces, no realtime scheduling, no personality (ABI) changes,
# no setuid/setgid binaries, no executable+writable memory.
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
# W^X (writable-XOR-executable memory). Rendered by deploy.sh from
# SUPERMUX_HARDENED (default 0 → `no`; set 1 → `yes`):
#   - no  (DEFAULT, dev-friendly): agents CAN run V8 JIT + WebAssembly build
#         tools — npm, React, Vite, Next, Jest, esbuild, etc. supermux exists to
#         host parallel CODING agents, and W^X=yes makes the kernel reject the
#         JIT's mprotect(PROT_EXEC) (ENOMEM), so every such build fails for the
#         service AND its children (inherited via NoNewPrivileges). This is the
#         ONE directive incompatible with the product's core use, so it's off by
#         default. Every OTHER sandbox directive above/below stays ON either way
#         (non-root, ProtectSystem=strict, ReadWritePaths jail, PrivateTmp,
#         dropped capabilities, NoNewPrivileges, syscall filter, …) — that is the
#         real sandbox; this knob alone does not meaningfully widen it.
#   - yes (SUPERMUX_HARDENED=1, strict): classic W^X for the security-paranoid;
#         JS/WASM builds will fail. Only choose this if you do NOT need agents to
#         run build tools.
MemoryDenyWriteExecute=__MEMORY_DENY_WRITE_EXECUTE__
# Drop all Linux capabilities — the service binds an unprivileged port and
# needs none.
CapabilityBoundingSet=
AmbientCapabilities=
# Only the device nodes a normal userspace process needs (/dev/null, zero,
# random, …); no raw device access.
PrivateDevices=true
# Confine syscalls to the standard non-privileged set plus exec (needed to
# spawn tmux/git/claude/codex children).
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources @obsolete
SystemCallErrorNumber=EPERM
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
